# etc/upstart/airbnb_web_flask.conf
# upstart script for /etc/init/airbnb.conf
# sudo  start airbnb_6
# debug: $ init-checkconf airbnb_web_flask.conf

description "gunicorn application server running AirBnB_clone"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
setuid ubuntu
setgid www-data

script
	cd /home/ubuntu/github.com/johncoleman83/AirBnB
	exec gunicorn --bind 127.0.0.1:8002 wsgi.wsgi_web_flask:web_flask.app
end script
